# deploy/roles/deploy-code/templates/pokemaps.service.j2
[Unit]
Description=Pokemaps

[Service]
Type=simple
User=centos
WorkingDirectory={{ checkout_location }}
ExecStart={{ checkout_location }}/findEmAll.rb
Restart=on-abort
